# Dipti_Shevalekar_002245703_01/packer/systemdSetup.service
[Unit]
Description=VM setup with Image
ConditionPathExists=/opt/webappUnzipped/Dipti_Shevalekar_002245703_01
After=network.target

[Service]
Type=simple
User=csye6225
Group=csye6225
WorkingDirectory=/opt/webappUnzipped/Dipti_Shevalekar_002245703_01/
ExecStart=node server.js
Restart=on-failure
RestartSec=3
StandardOutput=append:/var/log/webapp_output.log
StandardError=append:/var/log/webapp_error.log
SyslogIdentifier=csye6225

[Install]
WantedBy=multi-user.target
WantedBy=cloud-init.target